# deploy/systemd/agent-gov-ci-status-relay.timer
[Unit]
Description=Poll terminal AgentGov governance CI results every 30 seconds

[Timer]
OnBootSec=1min
OnUnitInactiveSec=30s
AccuracySec=1s
Persistent=true
Unit=agent-gov-ci-status-relay.service

[Install]
WantedBy=timers.target
