# deploy/systemd/agent-gov-ci-status-relay.service
[Unit]
Description=AgentGov terminal CI status relay to Multica
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
TimeoutStartSec=5min
User=agent-gov-ci-relay
Group=agent-gov-ci-relay
WorkingDirectory=/opt/agent-gov-ci-status-relay/current
Environment=HOME=/var/lib/agent-gov-ci-status-relay
Environment=AGENT_GOV_RELAY_STATE_DIR=/var/lib/agent-gov-ci-status-relay/state
Environment=MULTICA_PROFILE=ci-status-relay
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=XDG_CACHE_HOME=/var/lib/agent-gov-ci-status-relay/.cache
EnvironmentFile=/etc/agent-gov-ci-status-relay/relay.env
LoadCredential=github_token:/etc/agent-gov-ci-status-relay/github_token
ExecStartPre=/usr/bin/test -f /opt/agent-gov-ci-status-relay/current/scripts/agent_gov_ci_status_relay.py
ExecStart=/usr/bin/python3 /opt/agent-gov-ci-status-relay/current/scripts/agent_gov_ci_status_relay.py poll
UMask=0077
NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictRealtime=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ReadOnlyPaths=/opt/agent-gov-ci-status-relay
ReadWritePaths=/var/lib/agent-gov-ci-status-relay/state
ReadWritePaths=/var/lib/agent-gov-ci-status-relay/.multica
ReadWritePaths=/var/lib/agent-gov-ci-status-relay/.cache

[Install]
WantedBy=multi-user.target
